# .config/autostart/appimagelauncherd.desktop
[Desktop Entry]
Exec=/usr/bin/appimagelauncherd
Name=appimagelauncherd
Type=Application
